# Production/services/access/access.door.mesh/src/files/etc/init.d/services/access/vng.iot.services.access.door.mesh.service
[Unit]
Description=Access door mesh
After=network.target

[Service]
ExecStartPre=/bin/sleep 15
Type=forking
ExecStart = /etc/init.d/vng.services.access.door.mesh start
ExecStop = /etc/init.d/vng.services.access.door.mesh stop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
