# Production/services/light/light.hub.group.mesh/src/files/etc/init.d/services/light/vng.iot.services.light.hub.group.mesh.service
[Unit]
Description=Light hub group mesh
After=network.target

[Service]
ExecStartPre=/bin/sleep 12
Type=forking
ExecStart = /etc/init.d/vng.services.light.hub.group.mesh start
ExecStop = /etc/init.d/vng.services.light.hub.group.mesh stop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
