# Production/services/device/device.access.reader.status.mesh/src/files/etc/init.d/services/device/vng.iot.services.device.access.reader.status.mesh.service
[Unit]
Description= Device access reader status mesh
After=network.target

[Service]
ExecStartPre=/bin/sleep 12
Type=forking
ExecStart = /etc/init.d/vng.services.device.access.reader.status.mesh start
ExecStop = /etc/init.d/vng.services.device.access.reader.status.mesh stop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
